# env-deployment/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# Nginx web server
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/nginx_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/nginx_stdout.log
stdout_logfile_maxbytes=10MB

# FastAPI backend server (binary executable)
[program:fastapi]
command=/usr/local/bin/agrotech-server
directory=/app/server
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/fastapi_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/fastapi_stdout.log
stdout_logfile_maxbytes=10MB
environment=OLLAMA_URL="%(ENV_OLLAMA_URL)s"

# Group to manage both services
[group:agrotech]
programs=nginx,fastapi
priority=999
